# minnehack.service
[Unit]
Description=Backend for Minnehack

[Service]
Type=simple
ExecStart=/home/minnehack/Opaque/target/release/minnehack
Restart=always
User=minnehack
Group=minnehack
WorkingDirectory=/home/minnehack/Opaque

[Install]
WantedBy=multi-user.target
